# config/supervisor/conf.d/dropbear.conf
[program:dropbear]
command=/docker/scripts/dropbear.sh
priority=10
autorestart=true
redirect_stderr=true
