# config/supervisor/conf.d/serf.conf
[program:serf]
command=/docker/scripts/serf.sh
priority=20
autorestart=true
stopsignal=INT
